# ansible/roles/backups/files/backup.service
[Unit]
Description=Backup

[Service]
Type=oneshot
ExecStart=/opt/backup.sh